# docker/mysql/root/etc/s6-overlay/s6-rc.d/mysql-config/up
#!/bin/execlineb -P

/bin/sh -c "/etc/s6-overlay/s6-rc.d/mysql-config/init"